# lazy.ansible/.manala/services/ssh-auth-bridge/run
#!/bin/sh

exec socat \
    UNIX-LISTEN:/var/run/ssh-auth-bridge.sock,fork,unlink-early,mode=777 \
    UNIX-CONNECT:/var/run/ssh-auth.sock
